# rootfs/etc/systemd/system/nftables.service.d/override.conf
[Service]
RemainAfterExit=yes
ExecReload=
ExecReload=/usr/sbin/nft -f /etc/nftables.conf
ExecStop=
ExecStop=/usr/sbin/nft flush ruleset

